# vmbuilder.conf.d/guestroot.host-vz/etc/rc.d/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

for ctid in 10{1..3}; do
  if [[ -d /vz/root/${ctid} ]]; then
    vzctl stop ${ctid}
    vzctl destroy ${ctid}
  fi
  /home/kemukins/vzctl.sh ${ctid}
done
